# services/graph-rewriting-service/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

; The Fastify app is the only process supervisord manages.
; Neo4j is started by the entrypoint script and runs as its own self-managed
; Java process — supervisord cannot monitor it via neo4j console since it was
; already daemonised by 'neo4j start'.
[program:app]
command=node /app/dist/index.js
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
